# deploy/cloudflared-doh.service
[Unit]
Description=cloudflared DNS-over-HTTPS proxy (spud-router)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/cloudflared-doh
ExecStart=/usr/local/bin/cloudflared proxy-dns --port 5053 --address 127.0.0.1 --upstream ${DOH_UPSTREAM}
Restart=on-failure
RestartSec=2
User=nobody
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
